# celery-start/conf/supervisor_celery_flower.ini
[program:start-celery-flower]
command=python manage.py celery flower
directory=/path/to/celery-start
environment=PATH="/path/to/virtualenv/bin"
stdout_logfile=/path/to/project/logs/celery.flower.log
stderr_logfile=/path/to/project/logs/celery.flower.log
autostart=true
autorestart=true
priority=995