# celery-start/conf/supervisor_celery_beat.ini
[program:start-celery-beat]
command=python manage.py celery beat -l INFO
directory=/path/to/celery-start
environment=PATH="/path/to/virtualenv/bin"
stdout_logfile=/path/to/project/logs/celery.beat.log
stderr_logfile=/path/to/project/logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998