# celery-start/conf/supervisor_celery_worker.ini
[program:start-celery-worker]
command=python manage.py celery worker -l INFO
directory=/path/to/celery-start
environment=PATH="/path/to/virtualenv/bin"
stdout_logfile=/path/to/project/logs/celery.worker.log
stderr_logfile=/path/to/project/logs/celery.worker.log
autostart=true
autorestart=true
priority=997    ; 优先级(worker先于flower启动)